# templates/lib/systemd/system/dns-stats-compactor.service.erb
[Unit]
Description=DNS-STATS recording daemon
Documentation=man:compactor(1)

[Service]
Type=simple
CPUAffinity=<%= @cpuset.join(',') %>
MemoryAccounting=true
MemoryHigh=<%= @memory_high %>G
MemoryMax=<%= @memory_max %>G
<%- if @compactor_options -%>
ExecStart=/usr/bin/compactor <%= @compactor_options %>
<%- else -%>
ExecStart=/usr/bin/compactor
<%- end -%>
ExecReload=/bin/kill -HUP $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
